# systemd/snapback.timer
[Unit]
Description=Backup btrfs subvolume snapshots daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
