# systemd/snapback.timer
[Unit]
Description=Snapshot & backup btrfs subvolumes daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
